# docker/go/service/goweb.service
[Unit]
Description=goweb
[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/usr/local/bin/main
[Install]
WantedBy=multi-user.target
